# src/server/settingsservice/settingsservice.service
# SPDX-FileCopyrightText: Copyright 2024 LG Electronics Inc.
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=Pullpiri Settings Service
Documentation=https://github.com/eclipse-pullpiri/pullpiri
After=network.target etcd.service
Wants=etcd.service

[Service]
Type=simple
User=pullpiri
Group=pullpiri
WorkingDirectory=/opt/pullpiri/settingsservice
ExecStart=/opt/pullpiri/settingsservice/bin/settingsservice \
  --config=/etc/pullpiri/settings.yaml \
  --etcd-endpoints=localhost:2379 \
  --bind-address=0.0.0.0 \
  --bind-port=8080 \
  --log-level=info
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=settingsservice

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/pullpiri

# Environment
Environment=RUST_LOG=settingsservice=info,warn
Environment=RUST_BACKTRACE=1

[Install]
WantedBy=multi-user.target